# release 1.0/startup-script.service
[Unit]
Description=Startup User Script
After=syslog.target

[Service]
ExecStart=/home/root/release/startup.sh

[Install]
WantedBy=multi-user.target
